# templates/rpi/phone_home.service
[Unit]
Description=mercury phone home
After=network.target
[Service]
User=mercury
ExecStart=/home/mercury/phone_home.sh
Restart=always
[Install]
WantedBy=multi-user.target
